# supervisord.conf
[supervisord]
nodaemon=true                   ; Run supervisord in the foreground
logfile=/dev/null               ; Log supervisord output to stdout/stderr
logfile_maxbytes=0

[program:celery]
command=/usr/local/bin/celery -A celery_app worker --loglevel=info --concurrency=2
directory=/app
user=appuser
autorestart=true
priority=10
environment=PYTHONUNBUFFERED=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gunicorn]
command=/usr/local/bin/gunicorn --bind 0.0.0.0:8000 --workers 1 --threads 4 --timeout 1800 --graceful-timeout 30 --keep-alive 300 --log-level info --worker-class uvicorn.workers.UvicornWorker server:app
directory=/app                  ; Working directory for the app
user=appuser                    ; Run gunicorn as our non-root user
autorestart=true
priority=20
environment=PYTHONUNBUFFERED=1  ; Ensure Python output is sent straight to logs
stdout_logfile=/dev/stdout      ; Redirect gunicorn stdout to container stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr      ; Redirect gunicorn stderr to container stderr
stderr_logfile_maxbytes=0